# i3/95-dock-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/su stathis -c 'DISPLAY=:0 /home/stathis/.i3/dock-setup.sh &'"
